# test.service
[Unit]
Description=slogr client

[Service]
User=root
Restart=always
RestartSec=120s
ExecStart= /bin/bash -c /home/ubuntu/client.sh

